# supervisor.conf
[program:uwsgi]
command = /opt/venv/bin/uwsgi --ini /srv/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
